# uwb-follow.service
[Unit]
Description=UWB Follow - 自主跟随系统
After=network.target roscore.service
Wants=network.target

[Service]
Type=simple
User=nvidia
Group=nvidia
Environment="HOME=/home/nvidia"
Environment="ROS_MASTER_URI=http://localhost:11311"
Environment="ROS_HOSTNAME=localhost"
ExecStart=/bin/bash -c "source /home/nvidia/sim2real_master-feature-master_and_slave/install/setup.bash && python3 -u /home/nvidia/moon/uwb_follow.py"
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
