# vision/zed-obstacle.service
[Unit]
Description=ZED Obstacle Perception - /moon/obstacle
After=network.target roscore.service
Wants=network.target

[Service]
Type=simple
User=nvidia
Group=nvidia
Environment="HOME=/home/nvidia"
Environment="ROS_MASTER_URI=http://localhost:11311"
Environment="ROS_HOSTNAME=localhost"
# Optional: Environment="SIM2REAL_WS=/path/to/your/sim2real"
ExecStart=/bin/bash -c "source /home/nvidia/moon/scripts/sim2real_env.sh && moon_source_sim2real && python3 -u /home/nvidia/moon/vision/zed_obstacle_node.py"
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
